# workers/crontab
*/1 * * * * node code/kylew/2013-11-web-historian/workers/htmlfetcher.js